# docker/supervisor-laravel.conf
[program:php-fpm]
command=php-fpm -F -y /usr/local/etc/php-fpm.conf
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
priority=10
stdout_logfile=/var/log/supervisor/php-fpm.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=root

[program:nginx]
command=nginx -g "daemon off;"
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
priority=20
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/nginx-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=root

[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=default,notifications,embeddings
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
priority=30
stdout_logfile=/var/log/supervisor/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/queue-worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=app

[program:queue-worker-long]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=1 --max-time=3600 --queue=long-running
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=600
priority=31
stdout_logfile=/var/log/supervisor/queue-worker-long.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/queue-worker-long-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=app

[program:scheduler]
command=php /var/www/html/artisan schedule:work
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
priority=40
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/scheduler-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=app
